# systemd/user/dark-theme.service
[Unit]
Description=Set dark theme

[Service]
Type=simple
ExecStart=%h/bin/dark

[Install]
WantedBy=default.target
