# systemd/user/light-theme.service
[Unit]
Description=Set light theme

[Service]
Type=simple
ExecStart=%h/bin/light

[Install]
WantedBy=default.target
